# crontab/crontab.l
0 * * * * /bin/bash /Users/jab/jab/crontab/bin/save_profile.sh
0 1 * * * /bin/bash /Users/jab/jab/crontab/bin/updatedb.sh
